# copy/opt/local/lib/svc/manifest/rundeck.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='network/rundeck' type='service' version='0'>
    <instance name='rundeck' enabled='false'>
      <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
        <service_fmri value='svc:/system/filesystem/local'/>
      </dependency>
      <dependency name='network-service' grouping='require_all' restart_on='none' type='service'>
        <service_fmri value='svc:/network/service'/>
      </dependency>
      <dependency name='name-services' grouping='require_all' restart_on='none' type='service'>
        <service_fmri value='svc:/milestone/name-services'/>
      </dependency>
      <dependency name='identity' grouping='optional_all' restart_on='none' type='service'>
        <service_fmri value='svc:/system/identity:domain'/>
      </dependency>
      <dependency name='system-log' grouping='optional_all' restart_on='none' type='service'>
        <service_fmri value='svc:/system/system-log'/>
      </dependency>
      <dependency name='autofs' grouping='optional_all' restart_on='none' type='service'>
        <service_fmri value='svc:/system/filesystem/autofs'/>
      </dependency>
      <method_context project=':default' resource_pool=':default' working_directory='/opt/rundeck'>
        <method_credential group='rundeck' limit_privileges=':default' privileges=':default' supp_groups=':default' user='rundeck'/>
        <method_environment>
          <envvar name='HOME' value='/opt/rundeck'/>
        </method_environment>
      </method_context>
      <exec_method name='start' type='method' exec='/opt/local/lib/svc/method/rundeck start' timeout_seconds='0'>
        <method_context/>
      </exec_method>
      <exec_method name='stop' type='method' exec=':kill' timeout_seconds='900'>
        <method_context/>
      </exec_method>
      <property_group name='startd' type='framework'>
        <propval name='ignore_error' type='astring' value='core,signal'/>
      </property_group>
      <template>
        <common_name>
          <loctext xml:lang='C'>Rundeck</loctext>
        </common_name>
        <documentation>
          <doc_link name='rundeck.org' uri='http://rundeck.org/' />
        </documentation>
      </template>
    </instance>
    <stability value='Unstable'/>
  </service>
</service_bundle>
